# ansible/roles/murano-agent/templates/systemd-murano-agent.service.j2
[Unit]
Description=OpenStack Murano Agent

[Service]
Type=simple
ExecStart=/usr/bin/muranoagent --config-file={{ murano_agent_config_file }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
