# ansible/roles/jupyterlab/templates/jupyter-lab.service.j2
[Unit]
Description=JupyterLab

[Service]
ExecStart={{ anaconda_link_dir }}/bin/jupyter-lab --no-browser --ip=127.0.0.1

[Install]
WantedBy=default.target
